# crontab
*/2 * * * *   cd /var/www/html && php artisan schedule:run >> /var/log/cron.log 2>&1
* * * * *  www-data echo "Hello, world!" >> /var/log/cron.log 2>&1
